# launchd/com.user.worklog.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.worklog.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/senjo.yoshiki/Work/senjoysk/worklog/dist/worklog-weekly</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/senjo.yoshiki/Work/senjoysk/worklog/logs/weekly_report.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/senjo.yoshiki/Work/senjoysk/worklog/logs/weekly_report.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>WORKLOG_ROOT</key>
        <string>/Users/senjo.yoshiki/Work/senjoysk/worklog</string>
    </dict>
</dict>
</plist>
